# gen-diagnostic.service
[Unit]
Description=Generator Input Diagnostic
After=multi-user.target
Wants=multi-user.target

[Service]
Type=oneshot
# Create marker file to prove service started (before Python runs)
ExecStartPre=/bin/sh -c 'echo "Service started at $(date)" > /var/log/gen_service_started.txt'
ExecStartPre=/bin/sleep 30
ExecStart=/usr/bin/python2 /usr/local/bin/diagnose_inputs.py
# Create marker file when service completes
ExecStartPost=/bin/sh -c 'echo "Service completed at $(date)" >> /var/log/gen_service_started.txt'
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
